# .agent/services/com.pai.memory-lifecycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.memory-lifecycle</string>

    <key>Description</key>
    <string>PAI Memory Lifecycle - Nightly memory consolidation and maintenance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/delphijc/.claude/bin/pai-memory-lifecycle</string>
        <string>nightly</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/delphijc</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/delphijc/Library/Logs/pai-memory-lifecycle.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/delphijc/Library/Logs/pai-memory-lifecycle.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/delphijc</string>
        <key>PATH</key>
        <string>/Users/delphijc/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PAI_DIR</key>
        <string>/Users/delphijc/.claude</string>
        <key>ENABLE_MEMORY_HOOKS</key>
        <string>true</string>
        <key>MEMORY_HTTP_PORT</key>
        <string>4242</string>
    </dict>
</dict>
</plist>
