# .agent/services/com.pai.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.heartbeat</string>

    <key>Description</key>
    <string>PAI Heartbeat - Cross-project repository health monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/delphijc/.claude/bin/pai-heartbeat</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/delphijc</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/delphijc/Library/Logs/pai-heartbeat.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/delphijc/Library/Logs/pai-heartbeat.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/delphijc</string>
        <key>PATH</key>
        <string>/Users/delphijc/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PAI_DIR</key>
        <string>/Users/delphijc/.claude</string>
    </dict>
</dict>
</plist>
