# .agent/services/com.pai.service-monitor-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.service-monitor-server</string>

    <key>Description</key>
    <string>PAI Service Monitor Server - Backend for service metrics collection</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/delphijc/.claude/bin/pai-service-monitor-server</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/delphijc/Projects/service-monitor-dashboard/apps/server</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/delphijc/Library/Logs/pai-service-monitor-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/delphijc/Library/Logs/pai-service-monitor-server.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/delphijc</string>
        <key>PATH</key>
        <string>/Users/delphijc/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SERVICES_DOMAIN</key>
        <string>galatasary.lan</string>
    </dict>
</dict>
</plist>
