# .agent/services/com.pai.cyber-alert-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.cyber-alert-server</string>

    <key>Description</key>
    <string>PAI Cyber Alert Server - Backend API (port 4200)</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/delphijc/.claude/bin/pai-cyber-alert-server</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/delphijc/Projects/cyber-alert-mgr</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/delphijc/Library/Logs/pai-cyber-alert-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/delphijc/Library/Logs/pai-cyber-alert-server.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/delphijc</string>
        <key>PATH</key>
        <string>/Users/delphijc/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PORT</key>
        <string>4200</string>
    </dict>
</dict>
</plist>
