# .agent/services/com.pai.observability-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.observability-server</string>

    <key>Description</key>
    <string>PAI Observability Server - Agent activity backend (port 5172)</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/delphijc/.claude/bin/pai-observability-server</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/delphijc/Projects/observability/apps/server</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/delphijc/Library/Logs/pai-observability-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/delphijc/Library/Logs/pai-observability-server.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/delphijc</string>
        <key>PATH</key>
        <string>/Users/delphijc/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PORT</key>
        <string>5172</string>
        <key>PAI_DIR</key>
        <string>/Users/delphijc/.claude</string>
    </dict>
</dict>
</plist>
