# .agent/services/markdown-editor.service
[Unit]
Description=PAI Markdown Editor (Web-based Viewer)
After=network-online.target
Wants=network-online.target
PartOf=pai-infrastructure.target

[Service]
Type=simple
WorkingDirectory=__HOME__/Projects/markdown-editor
ExecStart=__HOME__/.bun/bin/bun run server.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment="NODE_ENV=production"
Environment="PORT=4444"

[Install]
WantedBy=pai-infrastructure.target
